# infra/onion-proxy/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
loglevel=info

[program:tor]
command=/usr/bin/tor -f /etc/tor/torrc
autorestart=true
startretries=3
priority=10
stdout_logfile=/var/log/supervisor/tor.log
stderr_logfile=/var/log/supervisor/tor.err
user=tor

[program:onion-proxy]
command=/usr/local/bin/onion-proxy
autorestart=true
startretries=3
priority=20
stdout_logfile=/var/log/supervisor/proxy.log
stderr_logfile=/var/log/supervisor/proxy.err
user=proxy
; Wait a few seconds for tor to bootstrap before launching the proxy.
startsecs=8
